# supervisor/supervisord.conf
[supervisord]
logfile=/dev/null
pidfile=/home/user/supervisord.pid
nodaemon=true
strip_ansi=false
childlogdir=/home/user/logs

[program:syncd]
command=/bin/bash -lc '/home/user/.venv/bin/python -m sync'
directory=/home/user
autostart=true
autorestart=true
startsecs=1
priority=10
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:nginx]
command=/bin/bash -lc '/usr/local/openresty/bin/openresty -p /home/user/nginx -c /home/user/nginx/nginx.conf -g "daemon off;"'
directory=/home/user/nginx
autostart=true
autorestart=true
priority=5
stopsignal=TERM
stdout_logfile=/home/user/logs/nginx.supervisor.out
stderr_logfile=/home/user/logs/nginx.supervisor.err
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x800x24 -nolisten tcp
autostart=true
autorestart=true
startsecs=1
priority=8
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:napcat]
command=/bin/bash -lc '/home/user/scripts/wait-for-sync.sh && /home/user/scripts/run-napcat.sh'
environment=DISPLAY=":1",LIBGL_ALWAYS_SOFTWARE="1",HOME="/app",XDG_CONFIG_HOME="/app/.config",QT_OPENGL="software"
directory=/home/user
user=1000
autostart=true
autorestart=true
priority=20
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:sin-proxy]
command=/bin/bash -lc '/home/user/scripts/wait-for-sync.sh && /home/user/scripts/run-sin-proxy.sh'
directory=/home/user
user=1000
autostart=true
autorestart=false
startsecs=0
startretries=0
priority=25
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:maibot]
command=/bin/bash -lc '/home/user/scripts/wait-for-sync.sh && /home/user/.venv/bin/python bot.py'
directory=/home/user/MaiBot-Core
autostart=true
autorestart=true
priority=28
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:maim-bot-adapters]
command=/bin/bash -lc '/home/user/scripts/wait-for-sync.sh && /home/user/.venv/bin/python main.py'
directory=/home/user/MaiBot-Adapter
autostart=true
autorestart=true
priority=30
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:filebrowser]
command=/bin/bash -lc '/home/user/scripts/wait-for-sync.sh && /home/user/filebrowser --address 0.0.0.0 --port 8888 --root / --baseurl /filebrowser --database /home/user/filebrowser-data/filebrowser.db'
directory=/home/user
autostart=true
autorestart=true
priority=40
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:gotty]
command=/bin/bash -lc '/home/user/gotty --address 127.0.0.1 --port 8080 --path /t --credential ${GOTTY_USERNAME:-admin}:${GOTTY_PASSWORD:-adminadminadmin} -w bash'
directory=/home/user
autostart=true
autorestart=true
priority=45
stopsignal=TERM
stdout_logfile=/home/user/logs/gotty.supervisor.out
stderr_logfile=/home/user/logs/gotty.supervisor.err
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
